# atualiza-no-login.service
[Unit]
Description=Executa atualizacao automatica da sessao

[Service]
Type=oneshot
ExecStart=%h/.local/bin/atualiza-no-login.sh
